# cmd/alpamon/command/configs/alpamon-restart.service
[Unit]
Description= Restart Alpamon service

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart alpamon.service